# com.slacksor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.slacksor</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Replace with your Python path, e.g. /usr/bin/python3 or /opt/homebrew/bin/python3 -->
        <string>/usr/bin/python3</string>
        <string>src/slacksor.py</string>
        <string>serve</string>
    </array>
    <!-- Replace with absolute path to slacksor project directory -->
    <key>WorkingDirectory</key>
    <string>/path/to/slacksor</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/slacksor.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/slacksor.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- Load .env from WorkingDirectory; override here if needed -->
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
